# kousu-laptop-mode/src/usr/lib/systemd/system/power-saver-idle.service
[Unit]
After=power-profiles-daemon.service
Description=Reduce power-profiles-daemon

[Service]
Type=oneshot
ExecStart=/usr/bin/powerprofilesctl-respect-holds set power-saver
