# kousu-laptop-mode/src/usr/lib/systemd/system/power-saver-idle.timer
[Unit]
Description=Go to low power mode periodically (unless inhibited)

[Timer]
OnUnitActiveSec=15m

[Install]
WantedBy=timers.target
