# kousu-desktop-gnome/src/etc/profile.d/wayland.sh


# needed to get most KDE/QT apps to play nice under Wayland
export QT_QPA_PLATFORM=xcb
