# kousu-laptop-mode/src/etc/udev/rules.d/power-saver.rules
SUBSYSTEM=="power_supply", ENV{POWER_SUPPLY_TYPE}=="Mains", ENV{POWER_SUPPLY_ONLINE}=="0", RUN+="/usr/bin/powerprofilesctl-respect-holds set power-saver"
