# kousu-base/src/etc/systemd/system/reflector.service.d/etckeeper.conf
[Unit]
Before=etckeeper-reflector.service
Wants=etckeeper-reflector.service
