# kousu-device-nigiri/src/etc/profile.d/gpg-agent.sh
#!/bin/sh

# gpg-agent(1):
# >  You should always add the following lines to your .bashrc or whatever initialization file is used
#export GPG_TTY=$(tty)
