# kousu-desktop-i3/src/etc/profile.d/xterminal.sh

# this isn't already in i3-sensible-terminal
export TERMINAL=sakura
